# root/etc/s6-overlay/s6-rc.d/init-ani-ar-jelly/up
/etc/s6-overlay/s6-rc.d/init-ani-ar-jelly/run